# modules/gs-scripts/files/retrieve_public_keys_often.cron
*/5 * * * * root /usr/local/sbin/retrieve_public_keys
