# modules/epylog/files/epylog.cron
#!/bin/sh
# Run epylog daily.
#
# $Id$
#
/usr/sbin/epylog --cron